# crisp-anonymization-daemon.service
[Unit]
Description=CRISP Anonymization Daemon
Documentation=https://github.com/COS301-SE-2025/CRISP
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=crisp
Group=crisp
WorkingDirectory=/opt/crisp
ExecStart=/opt/crisp/venv/bin/python /opt/crisp/crisp_anonymization_daemon.py --config /opt/crisp/daemon_config.json --log-file /var/log/crisp/anonymization_daemon.log
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/crisp /var/run/crisp

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

# Environment
Environment=PYTHONPATH=/opt/crisp
Environment=REDIS_HOST=localhost
Environment=REDIS_PORT=6379

[Install]
WantedBy=multi-user.target